# Reddit/com.rosehomeslv.reddit-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.rosehomeslv.reddit-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ryanrose/Downloads/Claude/Reddit/reddit_pipeline.sh</string>
    </array>

    <!-- Run at 9:45 AM Monday through Friday -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>45</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>45</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>45</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>45</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>45</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/reddit_pipeline_launchagent.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/reddit_pipeline_launchagent.log</string>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run even if missed while asleep -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
